# isle/rootfs/var/run/systemd/generator/local-fs.target.requires/-.mount
# Automatically generated by systemd-fstab-generator

[Unit]
Documentation=man:fstab(5) man:systemd-fstab-generator(8)
SourcePath=/etc/fstab
Before=local-fs.target
After=blockdev@dev-disk-by\x2duuid-f66324e3\x2daaba\x2d4f02\x2d826f\x2ddaee188cd4ff.target

[Mount]
Where=/
What=/dev/disk/by-uuid/f66324e3-aaba-4f02-826f-daee188cd4ff
Type=ext4
Options=errors=remount-ro
