# isle/rootfs/var/run/systemd/generator/local-fs.target.requires/boot-efi.mount
# Automatically generated by systemd-fstab-generator

[Unit]
Documentation=man:fstab(5) man:systemd-fstab-generator(8)
SourcePath=/etc/fstab
Before=local-fs.target
Requires=systemd-fsck@dev-disk-by\x2duuid-04CE\x2d9284.service
After=systemd-fsck@dev-disk-by\x2duuid-04CE\x2d9284.service
After=blockdev@dev-disk-by\x2duuid-04CE\x2d9284.target

[Mount]
Where=/boot/efi
What=/dev/disk/by-uuid/04CE-9284
Type=vfat
Options=umask=0077
